# configs/flutter-check.timer
[Unit]
Description=Timer pour la vérification Flutter quotidienne

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target